# ops/systemd/polyedge-morning-report.service
[Unit]
Description=PolyEdge Morning Summary Generator
After=polyedge.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/PolyEdge
ExecStart=/usr/bin/env python3 scripts/generate_morning_report.py
