# ops/systemd/polyedge-data-backend-ireland.service
[Unit]
Description=PolyEdge Data Backend (Ireland Ingest)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/PolyEdge
EnvironmentFile=-/etc/polyedge/data-backend-ireland.env
ExecStart=/home/ubuntu/PolyEdge/target/release/polyedge_data_backend ireland-ingest --dataset-root ${POLYEDGE_DATA_ROOT} --symbols ${POLYEDGE_SYMBOLS} --timeframes ${POLYEDGE_TIMEFRAMES} --udp-bind ${POLYEDGE_IRELAND_UDP_BIND} --snapshot-ms ${POLYEDGE_SNAPSHOT_MS} --discovery-refresh-sec ${POLYEDGE_DISCOVERY_REFRESH_SEC}
Restart=always
RestartSec=2
LimitNOFILE=1048576
StandardOutput=append:/var/log/polyedge/data-backend-ireland.log
StandardError=append:/var/log/polyedge/data-backend-ireland.err.log

[Install]
WantedBy=multi-user.target
