# ops/systemd/polyedge-data-backend-tokyo.service
[Unit]
Description=PolyEdge Data Backend (Tokyo Collector)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/PolyEdge
EnvironmentFile=-/etc/polyedge/data-backend-tokyo.env
ExecStart=/home/ubuntu/PolyEdge/target/release/polyedge_data_backend tokyo-collector --dataset-root /var/lib/polyedge-data --symbols BTCUSDT --bind 0.0.0.0:0 --ireland-udp 10.0.3.123:9801
Restart=always
RestartSec=2
LimitNOFILE=1048576
StandardOutput=append:/var/log/polyedge/data-backend-tokyo.log
StandardError=append:/var/log/polyedge/data-backend-tokyo.err.log

[Install]
WantedBy=multi-user.target
