# ops/systemd/polyedge-tokyo-disk-guard.service
[Unit]
Description=PolyEdge Tokyo disk guard
After=network.target

[Service]
Type=oneshot
User=root
Environment=ROOT_USAGE_PCT=85
Environment=JOURNAL_KEEP=200M
Environment=REPO_DIR=/home/ubuntu/PolyEdge
Environment=RELEASE_ROOT=/home/ubuntu/polyedge-releases
ExecStart=/usr/bin/env bash /usr/local/bin/polyedge-tokyo-disk-guard.sh
