# ops/systemd/polyedge-data-backend-api.service
[Unit]
Description=PolyEdge Data Backend API (Ireland)
After=network-online.target polyedge-data-backend-ireland.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/PolyEdge
EnvironmentFile=-/etc/polyedge/data-backend-api.env
ExecStart=/home/ubuntu/PolyEdge/target/release/polyedge_data_backend api --dataset-root ${POLYEDGE_DATA_ROOT} --bind ${POLYEDGE_API_BIND}
Restart=always
RestartSec=2
LimitNOFILE=1048576
StandardOutput=append:/var/log/polyedge/data-backend-api.log
StandardError=append:/var/log/polyedge/data-backend-api.err.log

[Install]
WantedBy=multi-user.target
