# ops/systemd/polyedge-morning-report.timer
[Unit]
Description=Run PolyEdge Morning Report Daily

[Timer]
OnCalendar=*-*-* 07:00:00 UTC
Persistent=true
Unit=polyedge-morning-report.service

[Install]
WantedBy=timers.target
